# scripts/simadmin.service
[Unit]
Description=SimAdmin Backend Service
After=network.target ModemManager.service
Wants=ModemManager.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/simadmin
ExecStart=/opt/simadmin/simadmin
Restart=always
RestartSec=5
Environment=DBUS_SYSTEM_BUS_ADDRESS=unix:path=/var/run/dbus/system_bus_socket
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simadmin

[Install]
WantedBy=multi-user.target
